# roles/systemdtimer/templates/service.service.j2
[Unit]
Description="{{ systemdtimer_servicename }}"

[Service]
Type=oneshot
ExecStart="{{ systemdtimer_command }}"
{% if systemdtimer_user is defined %}
User={{ systemdtimer_user }}
{% endif %}
