# roles/systemdtimer/templates/service.service.j2
[Unit]
Description="{{ systemdtimer_servicename }}"

[Service]
Type=oneshot
ExecStart={{ systemdtimer_command }}
