# roles/systemdtimer/templates/timer.timer.j2
[Unit]
Description="{{ systemdtimer_servicename }}"

[Timer]
OnCalendar={{ systemdtimer_oncalendar }}
Unit={{ systemdtimer_servicename }}.service

[Install]
WantedBy=multi-user.target
